# myserver.conf
description "Form builder service"
author "Vadim Shandrikov"
  
start on (local-filesystems and net-device-up)
stop on runlevel [06]
  
respawn
exec python /home/ubuntu/form_builder/server.py